# data_test/rule-repository/00000000-0000-4000-8000-000000000013.json
{
    "rule_id": "00000000-0000-4000-8000-000000000013",
    "rule_name": "Registry Run Keys",
    "description": "Detects Registry Run Keys",
    "category": "Persistence",
    "severity": "high",
    "confidence": 80,
    "base_risk_score": 60,
    "version": "1",
    "mitre_mapping": [
        {
            "tactic": "Persistence",
            "technique_id": "T1547"
        }
    ],
    "conditions": {
        "all": [
            {
                "field": "process.command_line",
                "operator": "contains",
                "value": "CurrentVersion\\Run"
            }
        ]
    },
    "recommended_actions": [
        {
            "step": 1,
            "action": "Investigate immediately."
        }
    ]
}